# workers/agent-driver/com.zion.world-host.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.zion.world-host</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/opt/node@22/bin/node</string>
        <string>/Users/kodyw/Projects/Zion/workers/agent-driver/host.js</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>ZION_URL</key>
        <string>https://kody-w.github.io/zion/</string>
        <key>NODE_PATH</key>
        <string>/Users/kodyw/Projects/Zion/workers/agent-driver/node_modules</string>
    </dict>

    <key>WorkingDirectory</key>
    <string>/Users/kodyw/Projects/Zion/workers/agent-driver</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>30</integer>

    <key>StandardOutPath</key>
    <string>/tmp/zion-host.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/zion-host.err</string>
</dict>
</plist>
